# system/node.service
[Unit]
Description=EchoTrace Node Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/default/echotrace-node
WorkingDirectory=/opt/echotrace-node
ExecStart=/usr/bin/env python3 -m pi_nodes.node_service
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
